# launch_agents/com.swq.codex-feishu-control.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.swq.codex-feishu-control</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>/Users/swq/Documents/New project 2/codex_feishu_control/server.py</string>
    <string>--config</string>
    <string>/Users/swq/Documents/New project 2/codex_feishu_control/config.json</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/swq/Documents/New project 2/codex_feishu_control</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/swq/Library/Logs/codex-feishu-control.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/swq/Library/Logs/codex-feishu-control.err.log</string>
</dict>
</plist>
